# recipes-thermostat/thermostat/Thermostat/thermostat.service
[Unit]
Description=My personal thermostat

[Service]
ExecStart=/usr/bin/PersonalThermostat
Restart=always

[Install]
WantedBy=multi-user.target