# recipes-thermostat/thermostat/files/thermostat.service
[Unit]
Description=My personal thermostat

[Service]
ExecStart=/usr/bin/Thermostat
Restart=always

[Install]
WantedBy=multi-user.target